# config/systemd/user/battery-low-notify.service
[Unit]
Description=Battery low notification (mako sticky)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-low-notify
